# dot_config/systemd/user/atuin.service
[Unit]
Description=Sync atuin

[Service]
Type=oneshot
RemainAfterExit=true
ExecStop=/usr/bin/atuin sync

[Install]
WantedBy=default.target
