# dot_config/systemd/user/check_pkg_updates.timer
[Unit]
Description=Check updates for packages

[Timer]
OnBootSec=3min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
